# deploy/systemd/qwen-model-operations-report.service
[Unit]
Description=Qwen and ModelPort privacy-preserving operations report
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tiammomo/projects/infra/local-inference-stack
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/wsl/lib
ExecStart=/home/tiammomo/projects/infra/local-inference-stack/scripts/operations-report.sh --hours 24 --provider local_qwen --save --fail-on-alert
Nice=10

[Install]
WantedBy=default.target
