# deploy/systemd/qwen-model-runtime.service.in
[Unit]
Description=Reconcile the selected local inference runtime after login or boot
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@PROJECT_ROOT@
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/wsl/lib
Environment=QWEN_BOOT_PROFILE=latency
ExecStart=@PROJECT_ROOT@/scripts/runtime-reconcile.sh
RemainAfterExit=yes
TimeoutStartSec=10m

[Install]
WantedBy=default.target
