# deploy/systemd/qwen-model-operations-report.timer.in
[Unit]
Description=Generate the daily local inference operations report

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m
Unit=qwen-model-operations-report.service

[Install]
WantedBy=timers.target
